# platform/broadcom/sonic-platform-modules-cel/debian/platform-modules-ds4101.postinst
depmod -a
/usr/local/bin/platform_api_mgnt.sh install
systemctl enable pddf-platform-init.service
systemctl start pddf-platform-init.service
systemctl enable ds4101-pddf-platform-monitor.service
systemctl start ds4101-pddf-platform-monitor.service