# platform/broadcom/sonic-platform-modules-cel/debian/platform-modules-moonstone.prerm
# Disable and stop platform services
systemctl disable pddf-platform-post.service
systemctl stop pddf-platform-post.service
systemctl disable pddf-platform-init.service
systemctl stop pddf-platform-init.service
