# platform/broadcom/sonic-platform-modules-cel/debian/platform-modules-redstone.postinst
depmod -a
systemctl enable platform-modules-redstone.service
systemctl start platform-modules-redstone.service

#/usr/local/bin/platform_api_mgnt.sh install

