# platform/broadcom/sonic-platform-modules-inspur/debian/sonic-platform-inspur-sc5630el.postinst
#!/bin/sh
depmod -a
systemctl enable inspur-sc5630el-platform.service
systemctl start inspur-sc5630el-platform.service
