# platform/barefoot/sonic-platform-modules-asterfusion/cx532p-t/service/cx532p-monitor.service
[Unit]
Description=Setting up Asterfusion CX532P
After=sysinit.target
Requires=systemd-modules-load.service
After=systemd-modules-load.service
Before=sysvinit.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/cx532p_monitor.sh start

[Install]
WantedBy=multi-user.target
