# platform/marvell/sonic-platform-wistron/es2227_54ts_p/service/es2227_54ts_p-oob-led.service
[Unit]
Description=Platform es2227_54ts_p OOB LED Service
After=pmon.service es2227_54ts_p-init.service
DefaultDependencies=no

[Service]
ExecStartPre=/usr/bin/sleep 20
ExecStart=/usr/local/bin/es2227_54ts_p-oob-led.sh
KillSignal=SIGKILL
SuccessExitStatus=SIGKILL
Restart=on-failure
RestartSec=5s
#StandardOutput=tty

[Install]
WantedBy=multi-user.target
