# files/image_config/process-reboot-cause/post-process-reboot-cause.service
[Unit]
Description=Post reboot cause handling service
Requires=database.service
After=database.service

[Service]
Type=simple
ExecStart=/usr/bin/post-process-reboot-cause
