# platform/broadcom/sonic-platform-modules-cel/moonstone/systemd/pddf-platform-post.service
[Unit]
Description=platform initialization post service
After=pmon.service

[Service]
Type=idle
ExecStart=-/usr/local/bin/post_pddf_init.sh
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
