# platform/broadcom/sonic-platform-modules-mitac/debian/sonic-platform-mitac-ly1200-32x.postinst
depmod -a
systemctl enable sonic-platform-mitac-ly1200-32x.service
systemctl start sonic-platform-mitac-ly1200-32x.service
