# dockers/docker-fpm-quagga/collectors/bgp_collector.conf
[program:bgp_collector.py]
command=/usr/bin/python /usr/bin/bgp_collector.py
priority=1
autostart=true
autorestart=false
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
