# dockers/docker-fpm-quagga/collectors/quagga_watcher.conf
[program:quagga_watcher.py]
command=/usr/bin/python /usr/bin/quagga_watcher.py
priority=1
autostart=true
autorestart=false
startsecs=0
stdout_logfile=syslog
stderr_logfile=syslog
