# files/image_config/syspoll/syspoll.service
[Unit]
Description=System information poller daemon

[Service]
Type=simple
ExecStart=/usr/bin/syspoll

[Install]
WantedBy=multi-user.target
