# files/image_config/cron.d/core_cleanup
# Attempts to clean up core files every 15 minutes
*/15 * * * * root /usr/bin/core_cleanup > /dev/null 2>&1

