# platform/marvell-arm64/sonic-platform-rd98dx35xx-ext/rd98dx35xx_ext/service/rd98dx35xx-ext-init.service
[Unit]
Description=Platform rd98dx35xx-ext Service
After=sysinit.target
Before=pmon.target

[Service]
ExecStart=/usr/local/bin/rd98dx35xx-ext-init.sh
StandardOutput=tty

[Install]
WantedBy=multi-user.target
