# files/image_config/memory_statistics/memory_statistics.service
[Unit]
Description=Memory Statistics Daemon
After=network.target

[Service]
ExecStart=/usr/bin/memorystatsd
Restart=on-failure
RestartSec=5
EnvironmentFile=-/etc/memorystats.conf

[Install]
WantedBy=multi-user.target
