# platform/broadcom/sonic-platform-modules-fs/n8550_32c/service/n8550-32c-platform-monitor-psu.service
[Unit]
Description=FS N8550-32C Platform Monitoring PSU service
Before=pmon.service
After=sysinit.target
DefaultDependencies=no

[Service]
ExecStart=/usr/local/bin/fs_n8550_32c_monitor_psu.py
KillSignal=SIGKILL
SuccessExitStatus=SIGKILL

# Resource Limitations
LimitCORE=infinity

[Install]
WantedBy=multi-user.target
